# deploy/kavkazkamen.service
[Unit]
Description=KavkazKamen Django (gunicorn)
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/kavkazkamen/backend
EnvironmentFile=/var/www/kavkazkamen/backend/.env
ExecStart=/var/www/kavkazkamen/backend/.venv/bin/gunicorn \
    --workers 3 \
    --bind 127.0.0.1:8000 \
    --timeout 60 \
    --access-logfile - \
    --error-logfile - \
    kavkazkamen.wsgi:application
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
